# deploy/systemd/socmint-worker.service
[Unit]
Description=SOCMINT queued scan worker
After=network.target socmint.service
Wants=socmint.service

[Service]
Type=oneshot
User=socmint
Group=socmint
WorkingDirectory=/opt/socmint
EnvironmentFile=/etc/socmint/socmint.env
Environment=SOCMINT_WORKER_MAX_JOBS=1
ExecStart=/opt/socmint/venv/bin/python -m src.socmint.main process-jobs --max-jobs=${SOCMINT_WORKER_MAX_JOBS}
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/socmint /var/log/socmint
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
